# nuxeo-opensocial/nuxeo-distribution-opensocial-jetty/src/main/scripts/postinstall
# First install 
if [ "$1" -eq "1" ]; then
  # Auto start at server boot
  chkconfig --add lmportal
fi
FQDN=`hostname -f`

#Stops lm portal
service lmportal stop > /dev/null 2>&1

#Delete Bundle cache file
BUNDLES_CACHES=$RPM_INSTALL_PREFIX/lm-portal/data/bundles.caches
[ -f $BUNDLES_CACHES ] && rm -f $BUNDLES_CACHES

MODULE_DIR=$RPM_INSTALL_PREFIX/web/modules
[ -d $MODULE_DIR ] && rm -rf $RPM_INSTALL_PREFIX/web/modules

#Configure host properties
sed -i -e "s/localhost/${FQDN}/" $RPM_INSTALL_PREFIX/lm-portal/config/portal.properties
sed -i -e "s/localhost:8080/${FQDN}:8080/" $RPM_INSTALL_PREFIX/lm-portal/config/opensocial.properties

service lmportal start